# files/etc/s6/syslogd/run
#!/bin/sh

set -eu

exec /sbin/syslogd -nS -O -
